# crontab
*/5 * * * * root /usr/src/twreporter/twreporter_cron.sh >> /var/log/cron.log 2>&1

