# crontab
*/5 * * * * root ./twreporter_cron >> /var/log/cron.log 2>&1

